# photo_frame_nas_sync.service
[Unit]
Description=Sync photos from NAS
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=robert
ExecStart=/home/robert/digital_photo_frame/sync_from_nas.sh
StandardOutput=journal
StandardError=journal
