# photo_frame_server.service
[Unit]
Description=Digital Photo Frame Viewer Server
After=network.target

[Service]
Type=simple
User=robert
WorkingDirectory=/home/robert/digital_photo_frame
Environment="PORT=8080"
ExecStart=/home/robert/digital_photo_frame/venv/bin/python /home/robert/digital_photo_frame/viewer_server.py /home/robert/digital_photo_frame/config_frame.yaml
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
