# photo_frame_nas_sync.timer
[Unit]
Description=Weekly photo sync from NAS (1h after NAS prepares at 03:00)

[Timer]
OnCalendar=Mon 04:00
Persistent=true

[Install]
WantedBy=timers.target
