# systemd/photo_frame_update.timer
[Unit]
Description=Daily auto-update for photo frame

[Timer]
OnCalendar=*-*-* 04:00:00
RandomizedDelaySec=1800
Persistent=true

[Install]
WantedBy=timers.target
